# pkg/DEBIAN/postinst
#!/bin/bash
set -e

. /usr/share/debconf/confmodule

INSTALL_DIR="/usr/local/bin"
CLI_SCRIPT="${INSTALL_DIR}/sshoyu_cli.sh"
CONFIG_DIR="/etc/sshoyu"
CONFIG_FILE="${CONFIG_DIR}/sshoyu.conf"
SUDOERS_FILE="/etc/sudoers.d/sshoyu"

# --- 1. Ler respostas do debconf ---
db_get sshoyu-server/ssh_host
SSH_HOST="$RET"

db_get sshoyu-server/ssh_port
SSH_PORT="$RET"

db_get sshoyu-server/ssh_user
SSH_USER="$RET"

if [ -z "$SSH_HOST" ]; then
    echo "ERRO: O hostname do servidor SSH não pode ser vazio." >&2
    exit 1
fi


# --- 2. Criar /etc/sshoyu/sshoyu.conf ---
mkdir -p "$CONFIG_DIR"
cat > "$CONFIG_FILE" <<EOF
SSH_HOST=${SSH_HOST}
SSH_PORT=${SSH_PORT}
SSH_USER=${SSH_USER}
CADDYFILE_PATH=/etc/caddy/Caddyfile
MIN_REMOTE_PORT=3000
LOCK_FILE_DIR=/tmp
EOF
chmod 644 "$CONFIG_FILE"
echo "✓ Configuração salva em ${CONFIG_FILE}"

# --- 3. Instalar os scripts ---
install -m 0755 /usr/share/sshoyu/sshoyu_cli.sh "$CLI_SCRIPT"
echo "✓ sshoyu_cli.sh instalado em ${CLI_SCRIPT}"

install -m 0755 /usr/share/sshoyu/sshoyu-monitor.sh /usr/local/bin/sshoyu-monitor.sh
echo "✓ sshoyu-monitor.sh instalado em /usr/local/bin/sshoyu-monitor.sh"

install -m 0755 /usr/share/sshoyu/sshoyu-admin.sh /usr/local/bin/sshoyu
echo "✓ sshoyu-admin instalado em /usr/local/bin/sshoyu"

# --- 4. Criar usuário do sistema ---
if ! id "$SSH_USER" &>/dev/null; then
    adduser --system --shell /bin/bash --home "/home/${SSH_USER}" --group "$SSH_USER"
    echo "✓ Usuário '${SSH_USER}' criado"
else
    echo "  Usuário '${SSH_USER}' já existe, pulando criação"
fi

# --- 5. Gerar par de chaves SSH ---
SSH_DIR="/home/${SSH_USER}/.ssh"
KEY_FILE="${SSH_DIR}/id_ed25519"

mkdir -p "$SSH_DIR"
chown "${SSH_USER}:${SSH_USER}" "$SSH_DIR"
chmod 700 "$SSH_DIR"

if [ ! -f "$KEY_FILE" ]; then
    sudo -u "$SSH_USER" ssh-keygen -t ed25519 -f "$KEY_FILE" -N "" -C "sshoyu@$(hostname)"
    echo "✓ Par de chaves SSH gerado em ${KEY_FILE}"
else
    echo "  Chave SSH já existe em ${KEY_FILE}, pulando geração"
fi

# --- 6. Configurar authorized_keys com forced command ---
AUTH_KEYS="${SSH_DIR}/authorized_keys"
FORCED_COMMAND="command=\"${CLI_SCRIPT}\",no-X11-forwarding,no-agent-forwarding,no-pty"
PUB_KEY=$(cat "${KEY_FILE}.pub")

# Evita duplicação: remove linha anterior do forced command se existir
if [ -f "$AUTH_KEYS" ]; then
    grep -v "sshoyu_cli.sh" "$AUTH_KEYS" > "${AUTH_KEYS}.tmp" && mv "${AUTH_KEYS}.tmp" "$AUTH_KEYS" || true
fi

echo "${FORCED_COMMAND} ${PUB_KEY}" >> "$AUTH_KEYS"
chown "${SSH_USER}:${SSH_USER}" "$AUTH_KEYS"
chmod 600 "$AUTH_KEYS"
echo "✓ authorized_keys configurado com forced command"

# --- 7. Instalar regras sudoers ---
cat > "$SUDOERS_FILE" <<EOF
# SShoyu: permite ao usuário ${SSH_USER} gerenciar o Caddyfile e recarregar o Caddy
${SSH_USER} ALL=(ALL) NOPASSWD: /usr/bin/tee /etc/caddy/Caddyfile
${SSH_USER} ALL=(ALL) NOPASSWD: /usr/bin/tee -a /etc/caddy/Caddyfile
${SSH_USER} ALL=(ALL) NOPASSWD: /bin/sed * /etc/caddy/Caddyfile
${SSH_USER} ALL=(ALL) NOPASSWD: /bin/mv /tmp/Caddyfile.tmp /etc/caddy/Caddyfile
${SSH_USER} ALL=(ALL) NOPASSWD: /bin/systemctl reload caddy
EOF
chmod 440 "$SUDOERS_FILE"
# Valida a sintaxe do sudoers antes de aplicar
if ! visudo -c -f "$SUDOERS_FILE" &>/dev/null; then
    echo "ERRO: Arquivo sudoers inválido. Removendo." >&2
    rm -f "$SUDOERS_FILE"
    exit 1
fi
echo "✓ Sudoers configurado em ${SUDOERS_FILE}"

# --- 8. Ativar serviço de monitoramento ---
systemctl daemon-reload
systemctl enable sshoyu-monitor
systemctl restart sshoyu-monitor
echo "✓ sshoyu-monitor ativado"

# --- 9. Gerar install.sh para os clientes ---
CADDY_WWW="/usr/share/caddy"
INSTALL_SCRIPT="${CADDY_WWW}/install.sh"
CLIENT_SCRIPT="${CADDY_WWW}/sshoyu"

# Gera o script sshoyu com as variáveis do servidor embutidas
sed \
    -e "s|__SSH_HOST__|${SSH_HOST}|g" \
    -e "s|__SSH_PORT__|${SSH_PORT}|g" \
    -e "s|__SSH_USER__|${SSH_USER}|g" \
    /usr/share/sshoyu/ssh_client.sh > "$CLIENT_SCRIPT"
chmod 644 "$CLIENT_SCRIPT"

# Gera o instalador curl-pipe
cat > "$INSTALL_SCRIPT" <<'INSTALL'
#!/bin/sh
set -e

INSTALL_DIR="/usr/local/bin"
KEY_DIR="$HOME/.ssh"
KEY_FILE="${KEY_DIR}/sshoyu_id_ed25519"

echo "[*] Installing sshoyu client..."

# Baixa o script sshoyu já configurado
curl -fsSL "http://__SSH_HOST__/sshoyu" | sudo tee "${INSTALL_DIR}/sshoyu" > /dev/null
sudo chmod +x "${INSTALL_DIR}/sshoyu"

# Gera chave SSH dedicada se não existir
if [ ! -f "$KEY_FILE" ]; then
    mkdir -p "$KEY_DIR"
    chmod 700 "$KEY_DIR"
    ssh-keygen -t ed25519 -f "$KEY_FILE" -N "" -C "sshoyu-client@$(hostname)"
    echo ""
    echo "[*] Chave SSH gerada. Envie a chave pública abaixo para o administrador"
    echo "    do servidor sshoyu para liberar seu acesso:"
    echo ""
    cat "${KEY_FILE}.pub"
    echo ""
else
    echo "[*] Chave SSH já existe em ${KEY_FILE}, pulando geração"
fi

echo "[+] sshoyu instalado em ${INSTALL_DIR}/sshoyu"
echo "[+] Use: sshoyu <subdomain> <localport>"
INSTALL

# Substitui o placeholder do host no install.sh
sed -i "s|__SSH_HOST__|${SSH_HOST}|g" "$INSTALL_SCRIPT"
chmod 644 "$INSTALL_SCRIPT"
echo "✓ install.sh gerado em ${INSTALL_SCRIPT}"
echo "✓ sshoyu client script gerado em ${CLIENT_SCRIPT}"

# Instala a página HTML de documentação
sed "s|abreudenis\.com|${SSH_HOST}|g" /usr/share/sshoyu/index.html > "${CADDY_WWW}/index.html"
chmod 644 "${CADDY_WWW}/index.html"
echo "✓ index.html gerado em ${CADDY_WWW}/index.html"

# --- 10. Resumo final ---
echo ""
echo "=============================================="
echo "  SShoyu Server instalado com sucesso!"
echo "=============================================="
echo ""
echo "  Servidor SSH: ${SSH_HOST}:${SSH_PORT}"
echo "  Usuário SSH:  ${SSH_USER}"
echo "  Config:       ${CONFIG_FILE}"
echo ""
echo "  Chave pública gerada (distribua para os clientes):"
echo ""
cat "${KEY_FILE}.pub"
echo ""
echo "  Os clientes devem configurar ssh_client.sh com:"
echo "    SSH_HOST=${SSH_HOST}"
echo "    SSH_PORT=${SSH_PORT}"
echo "    SSH_USER=${SSH_USER}"
echo "    SSH_KEY=<caminho para a chave privada do cliente>"
echo "=============================================="

db_stop
