# debian/prerm
#!/bin/bash
set -e

# Para o serviço de monitoramento
systemctl stop sshoyu-monitor 2>/dev/null || true
systemctl disable sshoyu-monitor 2>/dev/null || true

# Encerra túneis ativos
LOCK_DIR="/tmp"

for lock_file in "${LOCK_DIR}"/sshoyu_tunnel_*.lock; do
    [ -f "$lock_file" ] || continue
    pid=$(head -n1 "$lock_file" 2>/dev/null || true)
    if [ -n "$pid" ] && kill -0 "$pid" 2>/dev/null; then
        echo "Encerrando túnel ativo (PID ${pid})..."
        kill "$pid" 2>/dev/null || true
        sleep 1
    fi
    rm -f "$lock_file"
done
